# templates/default/puma.service.erb
# Generated by Chef. Changes will be overwritten.
#
[Unit]
Description=Puma HTTP Server for: <%= @app_dir %>
After=network.target

[Service]
Type=simple
User=<%= @user %>
WorkingDirectory=<%= @app_dir %>
EnvironmentFile=-<%= @env_file %>
ExecStart=/bin/bash -lc 'bundle exec --keep-file-descriptors puma -C <%= @conf_path %>'

Restart=always

[Install]
WantedBy=multi-user.target
